# inOrch-TMF-Proxy/systemd-portforward-inorch-tmf-proxy.service
#
# Systemd unit to keep a kubectl port-forward running so the inOrch-TMF-Proxy
# service inside Minikube is reachable from remote machines.
#
[Unit]
Description=inOrch-TMF-Proxy port-forward (expose svc/inorch-tmf-proxy on port 3020)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=telco
WorkingDirectory=/home/telco/arneme/INTEND-Project/5G4Data-public/inOrch-TMF-Proxy
Environment=KUBECONFIG=/home/telco/.kube/config
ExecStart=/usr/local/bin/kubectl -n inorch-tmf-proxy port-forward svc/inorch-tmf-proxy 3020:3020 --address 0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

